# crontab
*/5 * * * * root python /app/fetchData.py >> /app/etl.log 2>&1
