# crontab
* * * * * root python /app/fetchData.py >> /var/log/etl.log 2>&1
